# centos_supervisor_sshd/supervisord.conf
[supervisord]
nodaemon=True

[include]
files = /etc/supervisor/conf.d/*.conf

[program:sshd]
command=/sbin/sshd -D
